# linux-orchestrator/media-orchestrator.service
[Unit]
Description=Media Lifecycle Orchestrator Web UI & Pipeline Engine
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/media_orchestrator
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
